# data/lmkd-user-server.service
[Unit]
Description=Low Memory Killer user server

[Service]
Type=simple
ExecStart=/usr/libexec/lmkd-user-server
Restart=on-failure
